# rescue/oem/after/etc/skel/.config/autostart/readme.desktop
[Desktop Entry]
Type=Application
Version=1.0
Exec=/usr/bin/mousepad Desktop/readme.txt
Name=Open readme.txt