# rescue/oem/after/etc/skel/.config/autostart/readme.desktop
[Desktop Entry]
Type=Application
Version=1.0
Exec=mousepad Desktop/readme.txt
