# rescue/oem/after/etc/skel/.config/autostart/set-bg.desktop
[Desktop Entry]
Type=Application
Version=1.0
Exec=.config/autostart/set-bg.sh
Name=Set background