# rescue/oem/after/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \u' --autologin ubuntu --noclear %I $TERM
